# openwrt-contrib/etc/init.d/strongDNS2
#!/bin/sh /etc/rc.common

START=99
STOP=15
USE_PROCD=1

STRONGDNS2_BIN="/usr/bin/strongDNS2"
STRONGDNS2_DATA="/usr/share/strongDNS2"
MARK_SITES_NFT="$STRONGDNS2_DATA/hooks.d/mark_sites-nft.sh"

DEFAULT_QUEUENUM=1
DEFAULT_DEBUG=0
DEFAULT_USE_MARK_SITES=0
DEFAULT_DNS_SERVER_V4="{1.1.1.1,1.0.0.1,8.8.8.8,8.8.4.4}"
DEFAULT_DNS_SERVER_V6="{2606:4700:4700::1111,2606:4700:4700::1001}"

load_config() {
	. /lib/functions.sh

	config_load strongDNS2

	local enabled
	config_get enabled main enabled 1
	[ "$enabled" = "1" ] || return 1

	config_get QUEUENUM main queuenum "$DEFAULT_QUEUENUM"
	config_get USE_MARK_SITES main use_mark_sites "$DEFAULT_USE_MARK_SITES"
	config_get DNS_SERVER_V4 main dns_server_v4 "$DEFAULT_DNS_SERVER_V4"
	config_get DNS_SERVER_V6 main dns_server_v6 "$DEFAULT_DNS_SERVER_V6"
	config_get DEBUG main debug "$DEFAULT_DEBUG"

	return 0
}

start_service() {
	load_config || {
		echo "strongDNS2 disabled in /etc/config/strongDNS2"
		return 0
	}

	echo "Starting strongDNS2 service..."

	nft delete table inet strongDNS2 >/dev/null 2>&1
	cat <<EOF | nft -f -
define DNS_SERVER_V4 = $DNS_SERVER_V4
define DNS_SERVER_V6 = $DNS_SERVER_V6
table inet strongDNS2 {
  chain input {
    type filter hook input priority filter; policy accept;
    ip  saddr \$DNS_SERVER_V4 udp sport 53 counter queue to $QUEUENUM
    ip6 saddr \$DNS_SERVER_V6 udp sport 53 counter queue to $QUEUENUM
  }
  chain forward {
    type filter hook forward priority filter; policy accept;
    ip  saddr \$DNS_SERVER_V4 udp sport 53 counter queue to $QUEUENUM
    ip6 saddr \$DNS_SERVER_V6 udp sport 53 counter queue to $QUEUENUM
  }
}
EOF

	if [ "$USE_MARK_SITES" = "1" ] && [ -x "$MARK_SITES_NFT" ]; then
		"$MARK_SITES_NFT"
	fi

  local cmd_args
  cmd_args="-q $QUEUENUM"

  # debug
  [ "$DEBUG" = "1" ] && cmd_args="$cmd_args -d"

  if [ "$USE_MARK_SITES" = "1" ]; then
    cmd_args="$cmd_args -m"
  fi

	procd_open_instance
	procd_set_param command "$STRONGDNS2_BIN" $cmd_args
	procd_set_param respawn
	procd_set_param stdout 1
	procd_set_param stderr 1

	# 如果 ujail 和 capabilities 配置存在，则启用沙盒保护
	if [ -x /sbin/ujail ] && [ -e /etc/capabilities/strongDNS2.json ]; then
		procd_add_jail strongDNS2 ronly requirejail
		procd_set_param capabilities /etc/capabilities/strongDNS2.json
		procd_set_param user nobody
		procd_set_param group nogroup
		procd_set_param no_new_privs 1

		procd_add_jail_mount "$STRONGDNS2_DATA/mark_sites/douyin"
		procd_add_jail_mount "$STRONGDNS2_DATA/mark_sites/youtube"
		procd_add_jail_mount "$MARK_SITES_NFT"
		procd_add_jail_mount "$STRONGDNS2_DATA/ipv4.txt"
		procd_add_jail_mount "$STRONGDNS2_DATA/ipv6.txt"

		if command -v nft >/dev/null 2>&1; then
			procd_add_jail_mount "$(command -v nft)"
		fi
		procd_add_jail_mount "$STRONGDNS2_BIN"
	fi

	procd_close_instance
}

stop_service() {
	echo "Stopping strongDNS2 service..."
	nft delete table inet strongDNS2 >/dev/null 2>&1
}

reload_service() {
	stop_service
	start_service
}

# vim: set sw=2 ts=2 et:
